# deploy/docker-compose/systemd/greptile-saml.service
[Unit]
Description=Greptile SAML Service (Jackson)
Requires=greptile-app.service
After=greptile-app.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/greptile
ExecStart=/usr/bin/docker compose --profile saml up -d
ExecStop=/usr/bin/docker compose --profile saml down
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
