# terraform/stacks/aws-ec2/files/bootstrap/greptile-compose.service
[Unit]
Description=Greptile Docker Compose Stack
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/greptile
ExecStartPre=/opt/greptile/pull-secrets.sh
ExecStart=/usr/bin/docker compose --project-directory /opt/greptile up
ExecStop=/usr/bin/docker compose --project-directory /opt/greptile down
TimeoutStartSec=0
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
