# terraform/stacks/aws-ec2/files/bootstrap/greptile-compose.service
[Unit]
Description=Greptile Docker Compose Stack
Wants=network-online.target
After=network-online.target docker.service hatchet-token-setup.service
Requires=docker.service greptile-compose-hatchet.service
ConditionPathExists=/opt/greptile/.hatchet-token-ready

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/greptile
ExecStart=/usr/bin/docker compose --project-directory /opt/greptile --profile greptile up --detach
ExecStop=/usr/bin/docker compose --project-directory /opt/greptile --profile greptile stop
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
