# terraform/stacks/aws-ec2/files/bootstrap/hatchet-token-setup.service
[Unit]
Description=Generate Hatchet client token for Greptile
Wants=network-online.target
After=network-online.target docker.service greptile-compose-hatchet.service
Requires=docker.service greptile-compose-hatchet.service
ConditionPathExists=/opt/greptile/.env

[Service]
Type=oneshot
WorkingDirectory=/opt/greptile
ExecStart=/opt/greptile/bin/generate-hatchet-token.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
